# simetbox-openwrt-simet-nping/files/etc/uci-defaults/00-simet-nping-unprivileged
#!/bin/sh
chmod u+s /usr/bin/simet-nping || true
:
