# simetbox-openwrt-simet-ma/files/etc/uci-defaults/99-simet-ma
#!/bin/sh
# Copyright (C) 2018 NIC.br
# Released under the public domain or GPLv2+ at your option

crontab -l >/dev/null 2>&1 || exit 0
crontab -l 2>/dev/null | grep -q simetbox_register_ma.sh && exit 0

( crontab -l ;
  echo simetbox_register_ma.sh | awk '{ srand(); print int(60*rand()) " " int(24*rand()) " * * * " $0 }' ) | crontab -
:
