# simetbox-openwrt-simet-ma/files/etc/init.d/simet-ma
#!/bin/sh /etc/rc.common

# SIMET2-MA - SIMET2 Measurement Agent services
# Copyright (c) 2018 NIC.br
# Distributed under the GPLv2+

START=99
STOP=10
USE_PROCD=1

[ -r /usr/lib/simet/simet-ma.conf ] && . /usr/lib/simet/simet-ma.conf
[ -r /etc/simet/simet-ma.conf ] && . /etc/simet/simet-ma.conf
AGENT_ID_FILE=${AGENT_ID_FILE:-/etc/simet/agent-id}
AGENT_TOKEN_FILE=${AGENT_TOKEN_FILE:-/etc/simet/agent.jwt}
SIMET_INETUP_SERVER=${SIMET_INETUP_SERVER:-simet-monitor-inetup.simet.nic.br}

BOXID=$(get_mac_address.sh 2>/dev/null | tr A-F a-f | tr -d ': -') || true
BOOTID=$(cat /proc/sys/kernel/random/boot_id) || true

boot() {
	simetbox_register_ma.sh <&- 2>&- >&- &
	start "$@"
}

start_service() {
	procd_open_instance simet-inetuptime
	procd_set_param file "$AGENT_ID_FILE" "$AGENT_TOKEN_FILE" /usr/lib/simet/simet-ma.conf /etc/simet/simet-ma.conf
	procd_set_param user nobody
	procd_set_param pidfile /var/run/simet-inetupc.pid
	procd_set_param respawn 180 5 10
	procd_set_param stderr 1
	procd_set_param command /usr/bin/inetupc -M "${LMAP_TASK_NAME_PREFIX}inetupc"
	[ -r "$AGENT_TOKEN_FILE" ] && procd_append_param command -j "$(cat $AGENT_TOKEN_FILE)"
	[ -r "$AGENT_ID_FILE" ] && procd_append_param command -d "$(cat $AGENT_ID_FILE)"
	[ -n "$BOOTID" ] && procd_append_param command -b "$BOOTID"
	[ -n "$BOXID" ] && procd_append_param command -m "$BOXID"
	procd_append_param command "$SIMET_INETUP_SERVER"
	procd_close_instance
}

reload_service() {
	/etc/init.d/simet-lmapd reload || true
	restart "$@"
}

